# 04.update/to-2.5.7/from-2.3.x-and-2.4.x-and2.5.x/update/code/templete/svnadmind/svnadmind.service
[Unit]
Description=Subversion protocol daemon
After=syslog.target network.target

[Service]
Type=simple
ExecStart=/usr/bin/php /var/www/html/server/svnadmind.php start
ExecStop=/usr/bin/php /var/www/html/server/svnadmind.php stop

[Install]
WantedBy=multi-user.target