# 04.update/to-2.5.9/from-2.3.x-and-2.4.x-and2.5.x/update/code/templete/svnserve/svnserve.service
[Unit]
Description=Subversion protocol daemon
After=syslog.target network.target
 
[Service]
Type=forking
EnvironmentFile=%s
ExecStart=%s --daemon --pid-file=%s $OPTIONS
 
[Install]
WantedBy=multi-user.target